# systemd/sa-influx-backup.service
[Unit]
Description=Solar Assistant InfluxDB backup to FTP

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sa_influx_backup_to_ftp.sh
